# systemd/niri-wallpaper-restore.service
[Unit]
Description=Restore niri wallpaper state
After=awww-daemon.service
Requires=awww-daemon.service

[Service]
Type=oneshot
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
ExecCondition=/usr/bin/env sh -c '(command -v awww-daemon >/dev/null 2>&1 || command -v awww >/dev/null 2>&1) && [ -n "$WAYLAND_DISPLAY" ]'
ExecStart=%h/.config/niri/scripts/wallpaper-sync.sh restore

[Install]
WantedBy=default.target
